# config/autostart/com.github.hluk.copyq.desktop
[Desktop Entry]
Name=CopyQ
Exec=copyq --start-server
Type=Application
Terminal=false
X-GNOME-Autostart-Delay=3
